# root/etc/s6-overlay/s6-rc.d/rsyslog/run
#!/command/execlineb -P
# send stderr to stdout so s6 logs pick it up
fdmove -c 2 1

# foreground
s6-setuidgid root
/usr/sbin/rsyslogd -n

